# deploy/toolhub.service
[Unit]
Description=ToolHub Django Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/toolhub
Environment=DJANGO_SECRET_KEY=your-production-secret-key-here
Environment=DJANGO_DEBUG=False
Environment=DJANGO_ALLOWED_HOSTS=yourdomain.com,www.yourdomain.com
Environment=DJANGO_CSRF_TRUSTED_ORIGINS=https://yourdomain.com,https://www.yourdomain.com
Environment=DJANGO_CORS_ORIGINS=https://yourdomain.com,https://www.yourdomain.com
Environment=DATABASE_URL=postgres://toolhub_user:password@127.0.0.1:5432/toolhub
ExecStart=/var/www/toolhub/.venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:8000 --workers 4 --timeout 120
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target